# boot/service/setup-reboot.service
[Unit]
Description=Photobooth Setup Reboot
DefaultDependencies=no

[Service]
Type=simple
ExecStart=/boot/install/setup_rebooth.sh
